# .github/debian/amaru.service
[Unit]
Description=Amaru Cardano node
Documentation=https://github.com/pragma-org/amaru
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=amaru
Group=amaru
EnvironmentFile=-/etc/default/amaru
EnvironmentFile=-/etc/sysconfig/amaru
WorkingDirectory=/var/lib/amaru
StateDirectory=amaru
RuntimeDirectory=amaru
ExecStart=/usr/bin/amaru run
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
LimitNOFILE=65536
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
